# MyDNS/MyDNS-notifier.timer
[Unit]
Description=Timer for MyDNS-notifier
After=network.target

[Timer]
OnCalendar=*-*-* 00/3:00:00

[Install]
WantedBy=timers.target